# stage-hmi/06-comfile/files/backlight.service
[Unit]
Description=Turn on backlight for Comfile 15 inches

[Service]
Type=oneshot
ExecStart=/sbin/backlight.sh

[Install]
WantedBy=multi-user.target