# stage-base/14-rtc/files/config-rtc-mcp79410.service
[Unit]
Description=RTC configuration for MCP79410

[Service]
ExecStart=/sbin/configrtc.sh

[Install]
WantedBy=multi-user.target
