# stage-base/14-rtc/files/config-rtc-MCP79410.service
[Unit]
Description=RTC configuration for MCP79410

[Service]
ExecStart=/usr/sbin/configrtc.sh

[Install]
WantedBy=multi-user.target
